# deployment/systemd/core.service

[Unit]
Description=AgentOS Core
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/agentos/agentos-core
ExecStart=/usr/local/bin/uvicorn agentos_core.main:app --host 0.0.0.0 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
